# .crontab
*/30 * * * * /home/rafa/scripts/orderUGR.py >> /home/rafa/.cronlog/orderUGR.log
