# crontab
*/30 * * * * /home/orderUGR/pormiscojones.sh >> /home/rafa/.cronlog/orderUGR.log
