# crontab
*/30 * * * * /home/<username>/scripts/orderUGR.py >> /home/<username>/.cronlog/orderUGR.log
